# services/cb_parser.timer
[Unit]
Description=cryptobot

[Timer]
OnBootSec=15m
OnUnitActiveSec=5m
